# data/systemd/ados-uplink-router.service
[Unit]
Description=ADOS Uplink Router (ground-station profile)
After=NetworkManager.service
PartOf=ados-supervisor.service
Requires=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
# Cut over: the native uplink matrix is the DEFAULT. The install verifies the
# binary is present before this unit starts. The packaged Python uplink service
# stays on disk as an emergency fallback: create the fallback flag and the unit
# reverts to Python on restart. The native daemon also owns the AP, USB-gadget
# tether, and the ethernet/wifi/modem managers in-process, so its subsumed
# packaged units are masked while it is the default (handled by the installer).
# To revert the native matrix to Python for a bench:
#   touch /etc/ados/net-python-fallback && systemctl restart ados-uplink-router
# `exec` replaces the shell image so the CAP_NET_ADMIN/CAP_NET_RAW ambient caps
# below carry to the chosen binary unchanged.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/net-python-fallback ]; then exec /opt/ados/venv/bin/python -m ados.services.ground_station.uplink_router; else exec /opt/ados/bin/ados-net; fi'
Restart=on-failure
RestartSec=10
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Needs to read /proc/net and modify the routing table.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /var/lib/ados /run/ados /etc/ados

MemoryMax=96M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-uplink-router

[Install]
WantedBy=ados-supervisor.service
