# data/systemd/ados-wfb-rx.service
[Unit]
Description=ADOS WFB-ng RX (ground-station profile)
# After= the network and the mavlink unit, NOT the supervisor itself.
# An earlier revision had After=ados-supervisor.service which deadlocked
# the supervisor's own startup: when the supervisor's start_service
# called systemctl start during its init phase, systemd held the start
# request waiting for the supervisor to reach 'active' state — but the
# supervisor needed this child to start before it could call
# sd_notify(READY=1). PartOf= still binds the lifecycle so a supervisor
# stop tears the rx unit down with it.
After=network.target ados-mavlink.service
PartOf=ados-supervisor.service
ConditionKernelCommandLine=!no-wfb
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
User=root
# Flag-gated cutover: the native ground-station data-plane runs only when BOTH
# the operator enabled the flag file AND the binary is present; otherwise the
# packaged Python service is the default, so a routine upgrade never changes
# behavior. To turn the native receive plane on for a bench:
#   touch /etc/ados/groundlink-rust-enabled && systemctl restart ados-wfb-rx
# `exec` replaces the shell image so the CAP_NET_RAW/CAP_NET_ADMIN ambient caps
# below carry to the chosen binary unchanged.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/groundlink-rust-enabled ] && [ -x /opt/ados/bin/ados-groundlink ]; then exec /opt/ados/bin/ados-groundlink; else exec /opt/ados/venv/bin/python -m ados.services.ground_station.wfb_rx; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# WFB RX needs raw socket + monitor-mode WiFi adapter access
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

MemoryMax=64M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb-rx

[Install]
WantedBy=ados-supervisor.service
