# data/systemd/ados-rerun-sink.service
[Unit]
Description=ADOS Rerun Visualization Sink
After=ados-api.service
PartOf=ados-supervisor.service

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.rerun_sink
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados

MemoryMax=200M
CPUQuota=20%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-rerun-sink

[Install]
WantedBy=ados-supervisor.service
