# data/systemd/ados-mcp.service
[Unit]
Description=ADOS MCP Server
After=ados-api.service
PartOf=ados-supervisor.service

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.mcp
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /etc/ados/mcp

MemoryMax=128M
CPUQuota=50%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-mcp

[Install]
WantedBy=ados-supervisor.service
