# data/systemd/ados-oled.service
[Unit]
Description=ADOS OLED UI (ground-station profile)
PartOf=ados-supervisor.service
# Start only when a display surface was recognized at install/probe time.
# The marker is written by install-display-overlay.sh (bound SPI LCD, HDMI,
# or I2C OLED) and by ados-display-probe.service (confirmed panel), and is
# removed on the no-display path. Gating on it means the unit SKIPS CLEAN
# (condition unmet, not FAILED) on a board with no panel, instead of
# starting against absent hardware. The /dev/i2c-* + /dev/fb* glob it
# replaced false-positived on every board that merely exposes an I2C bus.
ConditionPathExists=/etc/ados/display.enabled
# Wait for the boot-time presence probe to settle its verdict (confirm or
# revert + marker state) before we read the marker.
After=ados-display-probe.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
# Flag-gated cutover: the native display writer runs only when BOTH the operator
# enabled the flag file AND the binary is present; otherwise the packaged Python
# service is the default, so a routine upgrade never changes behavior. To turn
# the native display writer on for a bench:
#   touch /etc/ados/display-rust-enabled && systemctl restart ados-oled
ExecStart=/bin/sh -c 'if [ -e /etc/ados/display-rust-enabled ] && [ -x /opt/ados/bin/ados-display ]; then exec /opt/ados/bin/ados-display; else exec /opt/ados/venv/bin/python -m ados.services.ui.oled_service; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

MemoryMax=48M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-oled

[Install]
WantedBy=ados-supervisor.service
