# data/systemd/ados-assist.service
[Unit]
Description=ADOS Assist Diagnostics Service
After=ados-api.service ados-mavlink.service ados-health.service
PartOf=ados-supervisor.service

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.assist
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados

MemoryMax=60M
CPUQuota=15%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-assist

[Install]
WantedBy=ados-supervisor.service
