# data/systemd/ados-health.service
[Unit]
Description=ADOS Health Monitor
After=network-online.target
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.health
Restart=always
RestartSec=1
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados

MemoryMax=20M
CPUQuota=5%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-health

[Install]
WantedBy=ados-supervisor.service
