# data/systemd/ados-api.service
[Unit]
Description=ADOS REST API
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.api
Restart=always
RestartSec=2
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
# /etc/systemd/system is writable for plugin unit drops; the supervisor
# writes ados-plugin-<id>.service files there during install and removes
# them during uninstall. /var/lib/ados is writable for plugin state and
# the plugin-data unpack target.
ReadWritePaths=/var/ados /run/ados /etc/ados /etc/systemd/system /var/lib/ados

MemoryMax=128M
CPUQuota=25%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-api

[Install]
WantedBy=ados-supervisor.service
