# data/systemd/ados-memory.service
[Unit]
Description=ADOS World Model Memory Service
After=ados-api.service ados-vision.service
PartOf=ados-supervisor.service

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.memory
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados

MemoryMax=400M
CPUQuota=30%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-memory

[Install]
WantedBy=ados-supervisor.service
