# data/systemd/ados-wfb-relay.service
[Unit]
Description=ADOS WFB-ng relay (mesh relay role)
After=ados-batman.service
Wants=ados-batman.service
PartOf=ados-supervisor.service
# Only runs when the role sentinel reads `relay`. role_manager owns this file.
ConditionPathExists=/etc/ados/mesh/role
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
# The native relay loop is the DEFAULT. The install verifies the binary is
# present before this unit starts. The packaged Python relay stays on disk as
# an emergency fallback for this mesh role: create the fallback flag and the
# unit reverts to Python on restart. The native binary takes the role on its
# `--role` argument. (The direct-role receive plane is native-only; this marker
# governs only the mesh relay/receiver roles and is not toggled by `ados rust`.)
#   touch /etc/ados/groundlink-python-fallback && systemctl restart ados-wfb-relay
# `exec` replaces the shell image so the ambient caps below carry unchanged.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/groundlink-python-fallback ]; then exec /opt/ados/venv/bin/python -m ados.services.ground_station.wfb_relay; else exec /opt/ados/bin/ados-groundlink --role relay; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Relay shells out to wfb_rx -f which needs raw socket + monitor mode.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

MemoryMax=64M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb-relay

[Install]
WantedBy=ados-supervisor.service
