# data/systemd/ados-wfb-rx.service
[Unit]
Description=ADOS WFB-ng RX (ground-station profile)
After=ados-supervisor.service
PartOf=ados-supervisor.service
ConditionKernelCommandLine=!no-wfb
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.ground_station.wfb_rx
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# WFB RX needs raw socket + monitor-mode WiFi adapter access
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

MemoryMax=64M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb-rx

[Install]
WantedBy=ados-supervisor.service
