# data/systemd/ados-video.service
[Unit]
Description=ADOS Video Pipeline
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
# The native video orchestrator is OPT-IN until it is camera-bench-validated:
# run it only when BOTH the binary is present AND the operator enabled the flag
# file. Default (no flag) = the packaged Python service, so a routine upgrade
# never silently swaps the live video pipeline for an unproven binary. To A/B
# test the native binary on the bench:
#   touch /etc/ados/video-rust-enabled && systemctl restart ados-video
# `exec` replaces the shell image so SupplementaryGroups=video carries over.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/video-rust-enabled ] && [ -x /opt/ados/bin/ados-video ]; then exec /opt/ados/bin/ados-video; else exec /opt/ados/venv/bin/python -m ados.services.video; fi'
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Video needs camera device access
SupplementaryGroups=video
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /tmp /etc/ados

MemoryMax=512M
# Backstop the file-descriptor count. The real fix for the historical
# fd leak is upstream (camera discovery no longer hands a stale node to
# the encoder, and stderr draining is rate-limited), but a generous
# ceiling keeps a future leak from crashing the service on EMFILE.
LimitNOFILE=8192
# CPUQuota removed: 50% was throttling 43% of scheduling periods
# (nr_throttled=6255/14471), causing ffmpeg keyframe burst stalls
# that broke the RTSP TCP connection to mediamtx after ~24 minutes.
# The SBC has 4 cores / 15GB RAM — no need to throttle video encoding.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-video

[Install]
WantedBy=ados-supervisor.service
