# data/systemd/ados-kiosk.service
[Unit]
Description=ADOS HDMI Kiosk (ground-station profile)
# Order after the services that serve the cockpit on :8080 (the native control
# front + the FastAPI app it proxies). This is only an ordering HINT — systemd
# "started" is not "HTTP-ready" — so the service ALSO polls the target URL until
# it is served before launching the browser (a browser pointed at a not-yet-ready
# URL sticks on an error page).
After=ados-mediamtx-gs.service ados-control.service ados-api.service
PartOf=ados-supervisor.service
# NOT ordered After=graphical.target: this unit is WantedBy=multi-user.target,
# and graphical.target is itself After=multi-user.target, so ordering after
# graphical.target forms a cycle (multi-user wants this -> this after graphical
# -> graphical after multi-user) that systemd breaks by DELETING this unit's
# start job — the kiosk then never auto-starts on boot. Instead the service
# starts at multi-user.target and WAITS for the DRM display and the desktop
# session itself (so it works on both a desktop box and a headless appliance).
#
# No ConditionPathExists on the DRM node either: the KMS display device can be
# created a few seconds after this unit starts (a boot race a one-shot condition
# loses); the service waits for a display and exits cleanly if none appears.
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
# No DISPLAY is set here: the appliance (cage) path uses cage's own DRM backend
# and the service strips DISPLAY/WAYLAND_DISPLAY so cage never tries to nest
# under an absent X server; the windowed (live-desktop) path derives DISPLAY /
# WAYLAND_DISPLAY from the detected login session instead.
Environment=ADOS_CONFIG=/etc/ados/config.yaml
Environment=XDG_RUNTIME_DIR=/run/user/0
EnvironmentFile=-/etc/ados/env
ExecStart=/opt/ados/venv/bin/python -m ados.services.kiosk.kiosk_service
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-kiosk

[Install]
WantedBy=multi-user.target
