# data/systemd/ados-oled.service
[Unit]
Description=ADOS OLED UI (ground-station profile)
After=ados-supervisor.service
PartOf=ados-supervisor.service
ConditionPathExistsGlob=/dev/i2c-*
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.ui.oled_service
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

MemoryMax=48M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-oled

[Install]
WantedBy=ados-supervisor.service
