# data/systemd/ados-api.service
[Unit]
Description=ADOS REST API
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.api
Restart=always
RestartSec=2
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
# /etc is writable to support runtime re-pair via the upstream wfb-ng
# bind protocol: wfb_keygen writes /etc/gs.key and /etc/drone.key on the
# GS, and wfb_bind_server.sh copies /etc/drone.key, /etc/wifibroadcast.cfg,
# and /etc/bind.yaml into place on the drone during the bind handshake.
# Also covers plugin unit drops at /etc/systemd/system (supervisor writes
# ados-plugin-<id>.service files during install and removes them during
# uninstall) and runtime config at /etc/ados. /var/lib/ados holds plugin
# state and the plugin-data unpack target.
ReadWritePaths=/var/ados /run/ados /etc /var/lib/ados

MemoryMax=128M
CPUQuota=25%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-api

[Install]
WantedBy=ados-supervisor.service
