# data/systemd/ados-supervisor.service
[Unit]
Description=ADOS Drone Agent Supervisor
# Deliberately NOT After=/Wants=network-online.target.
# The radio link (WFB-ng on RTL8812EU monitor mode) and the entire
# video pipeline (camera -> encoder -> tee -> wfb_tx | wfb_rx ->
# ffmpeg sidecar -> mediamtx -> LCD) operate locally. They do not
# need DHCP, DNS, or a default route. Only the cloud-relay heartbeat
# child needs internet, and it retries indefinitely on its own.
# Gating the supervisor on network-online previously left the LCD
# blank for hours when the home Wi-Fi failed to associate at boot.
# Per Rule 26: "the install must be 100% AUTOMATIC" — an operator
# on a deserted island with both rigs powered must see video.
After=local-fs.target
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=notify
# Prefer the prebuilt native orchestrator when it is present on disk, and fall
# back to the packaged entry point otherwise. The prebuilt fetch is best-effort
# (offline / non-arm64 hosts skip it), so the unit must never point at a binary
# that may be absent. `exec` hands the process identity to the chosen binary so
# Type=notify readiness, the watchdog, and the cgroup all track the real PID.
ExecStart=/bin/sh -c 'if [ -x /opt/ados/bin/ados-supervisor ]; then exec /opt/ados/bin/ados-supervisor; else exec /opt/ados/venv/bin/ados-supervisor; fi'
Restart=always
RestartSec=1
WatchdogSec=30
TimeoutStartSec=180
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# /etc is writable so the radio bind handshake can write its key material:
# wfb_keygen produces /etc/gs.key + /etc/drone.key, and the bind-protocol
# helper scripts drop /etc/drone.key, /etc/wifibroadcast.cfg, and
# /etc/bind.yaml into place. Without it ProtectSystem=strict returns EROFS
# and the bind fails on keygen before any radio frame flows. PrivateTmp=yes
# (above) already gives the helper scripts a writable /tmp for their mktemp
# staging dir, so /tmp does not need to be listed here.
ReadWritePaths=/var/ados /var/lib/ados /etc/ados /run/ados /etc

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-supervisor

[Install]
WantedBy=multi-user.target
