# data/systemd/ados-usb-otg-host.service
[Unit]
Description=ADOS force USB OTG controllers into host mode
# Some single-board computers default their USB-C / OTG controller to the
# USB device (peripheral) role at every boot. In that role a connected
# powered hub's downstream peripherals never enumerate, so a USB WiFi radio
# adapter or a USB camera plugged through a hub is invisible to the agent.
# Writing usb_host to the controller's sysfs otg_role node forces host mode
# so the hub and its peripherals enumerate. The unit globs every matching
# otg_role node and no-ops on boards that expose none (the common case),
# so enabling it everywhere is safe.
DefaultDependencies=no
After=systemd-modules-load.service
Before=ados-supervisor.service network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=45
# Fast-bail unless this is an Allwinner SoC (per device-tree compatible), so
# the unit is an instant no-op on Rockchip / Pi / others even though it ships
# everywhere. On Allwinner the usbc driver settles the otg_role node a few
# seconds into boot and may set it to device mode, so retry for up to ~30 s:
# re-glob each node and re-assert usb_host until every present node reads host.
# A single write at a fixed boot point loses this race; the retry wins it.
ExecStart=/bin/sh -c 'grep -qiE "allwinner|sunxi" /proc/device-tree/compatible 2>/dev/null || exit 0; i=0; while [ "$i" -lt 30 ]; do pending=0; any=0; for n in /sys/bus/platform/devices/*.usbc*/otg_role /sys/devices/platform/soc*/*.usbc*/otg_role; do [ -e "$n" ] || continue; any=1; cur=$(cat "$n" 2>/dev/null || echo); [ "$cur" = usb_device ] && echo usb_host > "$n" 2>/dev/null; cur=$(cat "$n" 2>/dev/null || echo); [ "$cur" = usb_host ] || pending=1; done; { [ "$any" = 1 ] && [ "$pending" = 0 ]; } && exit 0; i=$((i + 1)); sleep 1; done; exit 0'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-usb-otg-host

[Install]
WantedBy=multi-user.target
