# data/systemd/ados-wfb-receiver.service
[Unit]
Description=ADOS WFB-ng receiver (mesh receiver role)
After=ados-batman.service
Wants=ados-batman.service
PartOf=ados-supervisor.service
# Only runs when the role sentinel reads `receiver`. role_manager owns this file.
ConditionPathExists=/etc/ados/mesh/role
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
# Flag-gated cutover, mirroring ados-wfb-rx.service: the native receiver loop
# runs only when BOTH the operator enabled the flag file AND the binary is
# present; otherwise the packaged Python receiver is the default, so a routine
# upgrade never changes behavior. The native binary takes the role on its
# `--role` argument. To turn the native receiver on for a bench:
#   touch /etc/ados/groundlink-rust-enabled && systemctl restart ados-wfb-receiver
# `exec` replaces the shell image so the ambient caps below carry unchanged.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/groundlink-rust-enabled ] && [ -x /opt/ados/bin/ados-groundlink ]; then exec /opt/ados/bin/ados-groundlink --role receiver; else exec /opt/ados/venv/bin/python -m ados.services.ground_station.wfb_receiver; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Receiver drives wfb_rx -a which binds a raw socket on the local NIC
# plus a UDP listener on the mesh interface.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

MemoryMax=96M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb-receiver

[Install]
WantedBy=ados-supervisor.service
