# data/systemd/ados-cloud.service
[Unit]
Description=ADOS Cloud Relay (MQTT + Convex Heartbeat)
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
# Flag-gated cutover: the native cloud relay runs only when BOTH the operator
# enabled the flag file AND the binary is present; otherwise the Python module
# is the default fallback, so a routine upgrade never changes behavior. To turn
# the native relay on for a bench:
#   touch /etc/ados/cloud-rust-enabled && systemctl restart ados-cloud
ExecStart=/bin/sh -c 'if [ -e /etc/ados/cloud-rust-enabled ] && [ -x /opt/ados/bin/ados-cloud ]; then exec /opt/ados/bin/ados-cloud; else exec /opt/ados/venv/bin/python -m ados.services.cloud; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
# Private /tmp for the bind orchestrator. wfb_bind_client.sh +
# wfb_bind_server.sh both use `mktemp -d` to stage the tar.gz of
# config + drone.key during the bind protocol. Without PrivateTmp,
# ProtectSystem=strict makes /tmp read-only and mktemp fails.
PrivateTmp=yes
# /etc is in the writable set because the WFB bind orchestrator hosted
# in this service spawns wfb_keygen + the upstream bind-protocol scripts
# which write /etc/drone.key, /etc/gs.key, /etc/wifibroadcast.cfg, and
# /etc/bind.yaml. Loosens the sandbox during the (rare) bind window.
ReadWritePaths=/var/ados /var/lib/ados /run/ados /etc/ados /etc

# The bind orchestrator calls `systemctl stop ados-wfb` /
# `systemctl start wifibroadcast@gs_bind` to flip wfb-ng profiles. Those
# hops require CAP_SYS_ADMIN via PID 1; bumping the bounding set is
# the systemd-friendly way to permit it without dropping ProtectSystem.
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN

# Memory cap sized for the bind orchestrator's worst case: the cloud
# service itself (~30M Python runtime + plugins) plus the bind window's
# transient children (wfb bind server ~30M, socat + wfb_tx + wfb_rx
# ~10M each). Anything under ~100M caused the bind to silently OOM-thrash
# inside the cgroup — wfb-server stuck in mem_cgroup swap reclaim, never
# responding to the VERSION/BIND handshake. 256M leaves comfortable
# headroom and matches the sibling video / mediamtx-gs scale.
MemoryMax=256M
# CPU cap removed — the bind protocol is bursty (a few seconds of work
# every 60 s when unpaired, near-zero post-pair), so capping was
# producing throttling exactly when the bind needed CPU to drive the
# bind protocol's base64 + tar + sha1sum + key write through socat.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-cloud

[Install]
WantedBy=ados-supervisor.service
