# data/systemd/ados-wfb.service
[Unit]
Description=ADOS WFB-ng Link Manager
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
# The native radio manager is OPT-IN until it reaches feature parity with the
# Python service: run it only when BOTH the binary is present AND the operator
# has explicitly enabled it via the flag file. Default (no flag) = the packaged
# Python service, so a routine upgrade never silently swaps the live radio for
# an under-parity binary. To A/B test the native binary on the bench:
#   touch /etc/ados/wfb-rust-enabled && systemctl restart ados-wfb
# `exec` replaces the shell image in place so the CAP_NET_RAW/CAP_NET_ADMIN
# ambient caps below carry to the chosen binary unchanged.
ExecStart=/bin/sh -c 'if [ -e /etc/ados/wfb-rust-enabled ] && [ -x /opt/ados/bin/ados-radio ]; then exec /opt/ados/bin/ados-radio; else exec /opt/ados/venv/bin/python -m ados.services.wfb; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# WFB needs raw socket + WiFi adapter access
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

# The original 30M cap was sized for the bare manager + wfb_tx subprocess.
# Adding the hop supervisor (band scanner, control-port broadcaster) plus
# the bitrate controller pushes steady-state memory past 29M, which puts
# the cgroup into hard reclaim and stalls the asyncio loop — the periodic
# hop timer effectively stops firing because allocations queue. 96M gives
# the supervisor's iw-scan output (~one channel-survey buffer) plus the
# control-port broadcaster headroom without inviting unbounded growth.
MemoryMax=96M
CPUQuota=25%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb

[Install]
WantedBy=ados-supervisor.service
