# data/systemd/ados-modem.service
[Unit]
Description=ADOS Modem Manager (ground-station profile)
After=ados-supervisor.service ModemManager.service
PartOf=ados-supervisor.service
# Only start when a USB modem is present. SIM7600G-H enumerates as
# /dev/ttyUSB0..3 and wwan0 on MBIM/QMI. The pattern keeps the unit
# inactive on ground stations with no cellular hardware.
ConditionPathExistsGlob=/dev/ttyUSB*

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.ground_station.modem_manager
Restart=on-failure
RestartSec=10
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /etc/ados /dev

MemoryMax=96M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-modem

[Install]
WantedBy=ados-supervisor.service
