# data/systemd/ados-video.service
[Unit]
Description=ADOS Video Pipeline
After=ados-mavlink.service
PartOf=ados-supervisor.service

[Service]
Type=simple
ExecStart=/opt/ados/venv/bin/python -m ados.services.video
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Video needs camera device access
SupplementaryGroups=video
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /tmp

MemoryMax=512M
# CPUQuota removed: 50% was throttling 43% of scheduling periods
# (nr_throttled=6255/14471), causing ffmpeg keyframe burst stalls
# that broke the RTSP TCP connection to mediamtx after ~24 minutes.
# The SBC has 4 cores / 15GB RAM — no need to throttle video encoding.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-video

[Install]
WantedBy=ados-supervisor.service
