# data/systemd/ados-video.service
[Unit]
Description=ADOS Video Pipeline
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
# The native video orchestrator is the only implementation; the standalone
# Python video service was removed once the binary was bench-validated. exec the
# binary unconditionally (failing loudly if it is absent rather than running a
# stale path); the installer always provisions it. `exec` replaces the shell
# image so SupplementaryGroups=video carries over. The reusable video library
# modules (pipeline, encoder, mediamtx, camera_mgr, local_tap, sei_*) stay in
# Python for the in-process demo pipeline and the ground-station consumers.
ExecStart=/bin/sh -c 'if [ -x /opt/ados/bin/ados-video ]; then exec /opt/ados/bin/ados-video; else echo "ados-video binary missing; reinstall to fetch it" >&2; exit 1; fi'
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# Video needs camera device access
SupplementaryGroups=video
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /tmp /etc/ados

MemoryMax=512M
# Backstop the file-descriptor count. The real fix for the historical
# fd leak is upstream (camera discovery no longer hands a stale node to
# the encoder, and stderr draining is rate-limited), but a generous
# ceiling keeps a future leak from crashing the service on EMFILE.
LimitNOFILE=8192
# CPUQuota removed: 50% was throttling 43% of scheduling periods
# (nr_throttled=6255/14471), causing ffmpeg keyframe burst stalls
# that broke the RTSP TCP connection to mediamtx after ~24 minutes.
# The SBC has 4 cores / 15GB RAM — no need to throttle video encoding.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-video

[Install]
WantedBy=ados-supervisor.service
