# data/systemd/ados-wfb.service
[Unit]
Description=ADOS WFB-ng Link Manager
After=ados-mavlink.service
PartOf=ados-supervisor.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
# Prefer the prebuilt native radio manager when present, fall back to the
# packaged Python service otherwise (the prebuilt fetch is best-effort). `exec`
# replaces the shell image in place so the CAP_NET_RAW/CAP_NET_ADMIN ambient
# caps below carry to the chosen binary unchanged.
ExecStart=/bin/sh -c 'if [ -x /opt/ados/bin/ados-radio ]; then exec /opt/ados/bin/ados-radio; else exec /opt/ados/venv/bin/python -m ados.services.wfb; fi'
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

# WFB needs raw socket + WiFi adapter access
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /dev /etc/ados

# The original 30M cap was sized for the bare manager + wfb_tx subprocess.
# Adding the hop supervisor (band scanner, control-port broadcaster) plus
# the bitrate controller pushes steady-state memory past 29M, which puts
# the cgroup into hard reclaim and stalls the asyncio loop — the periodic
# hop timer effectively stops firing because allocations queue. 96M gives
# the supervisor's iw-scan output (~one channel-survey buffer) plus the
# control-port broadcaster headroom without inviting unbounded growth.
MemoryMax=96M
CPUQuota=25%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-wfb

[Install]
WantedBy=ados-supervisor.service
