# data/systemd/ados-pic.service
[Unit]
Description=ADOS Pilot-in-Command Arbiter (ground-station profile)
After=ados-input.service
PartOf=ados-supervisor.service
# Skip cleanly when the binary was not fetched for this profile rather than
# crash-looping; the installer health gate is the loud failure for a binary
# that should be present on a ground station.
ConditionPathExists=/opt/ados/bin/ados-pic
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
User=root
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env
# Native-only: the packaged Python arbiter was deleted, so the daemon runs its
# binary directly. It also reads the front-panel GPIO buttons in-process, so the
# installer masks the packaged ados-buttons unit on a ground station.
ExecStart=/opt/ados/bin/ados-pic
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-pic

[Install]
WantedBy=multi-user.target
