# data/systemd/ados-oled.service
[Unit]
Description=ADOS OLED UI (ground-station profile)
After=ados-supervisor.service
PartOf=ados-supervisor.service
# Start when EITHER an I2C bus exists (small OLED carousel) OR a SPI
# framebuffer is bound (LCD page mode). systemd treats multiple
# ConditionPathExistsGlob lines as OR-of-any-match, so listing both
# means the service starts on a board with EITHER surface present.
ConditionPathExistsGlob=|/dev/i2c-*
ConditionPathExistsGlob=|/dev/fb*
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.ui.oled_service
Restart=on-failure
RestartSec=5
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

MemoryMax=48M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-oled

[Install]
WantedBy=ados-supervisor.service
