# data/systemd/ados-vision.service
[Unit]
Description=ADOS Vision Engine
After=ados-mavlink.service
PartOf=ados-supervisor.service
# Skip cleanly when the binary was not fetched for this profile rather than
# crash-looping "binary missing"; the installer's health gate is the loud
# failure for a binary that SHOULD be present.
ConditionPathExists=/opt/ados/bin/ados-vision
StartLimitInterval=60s
StartLimitBurst=5

[Service]
# Join the shared accounting slice so MemoryCurrent/CPUUsage are
# observable per unit (the slice sets MemoryAccounting=yes).
Slice=ados.slice
Type=simple
# The native vision engine is the only implementation; there is no Python
# fallback. exec the binary when present (failing loudly if it is absent rather
# than running a stale path); the installer provisions it best-effort. `exec`
# replaces the shell image so the chosen binary owns the cgroup + restart track.
ExecStart=/bin/sh -c 'if [ -x /opt/ados/bin/ados-vision ]; then exec /opt/ados/bin/ados-vision; else echo "ados-vision binary missing; reinstall to fetch it" >&2; exit 1; fi'
Restart=on-failure
RestartSec=3
Environment=ADOS_CONFIG=/etc/ados/config.yaml
# The onnx-enabled build loads the ONNX Runtime dynamically; point it at the
# shared library the installer places (absent on the default musl build, which
# ignores this). A stale path is harmless — the default build never dlopens ORT.
Environment=ORT_DYLIB_PATH=/opt/ados/lib/libonnxruntime.so
EnvironmentFile=-/etc/ados/env

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /etc/ados /dev/shm

MemoryMax=256M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-vision

[Install]
WantedBy=ados-supervisor.service
