# data/systemd/ados-ethernet.service
[Unit]
Description=ADOS Ethernet Manager (ground-station profile)
After=ados-supervisor.service NetworkManager.service
PartOf=ados-supervisor.service
Requires=NetworkManager.service
StartLimitInterval=60s
StartLimitBurst=5

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.ground_station.ethernet_manager
Restart=on-failure
RestartSec=10
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/ados /run/ados /etc/ados /var/lock

MemoryMax=64M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-ethernet

[Install]
WantedBy=ados-supervisor.service
