# data/systemd/ados-setup-captive.service
[Unit]
Description=ADOS Captive Portal DNS + HTTP responder (first-boot only)
After=ados-hostapd.service
PartOf=ados-supervisor.service
# Only run while setup is incomplete. Once the sentinel exists the unit
# is skipped (shown inactive, condition unmet) instead of reporting a
# spurious failed state on every already-configured boot.
ConditionPathExists=!/var/lib/ados/setup-complete

[Service]
Type=simple
User=root
ExecStart=/opt/ados/venv/bin/python -m ados.services.setup_webapp.captive_dns
Restart=no
Environment=ADOS_CONFIG=/etc/ados/config.yaml
EnvironmentFile=-/etc/ados/env

AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=no

MemoryMax=48M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ados-setup-captive

[Install]
WantedBy=ados-supervisor.service
