# agents/lite-rs/imagebuilder/overlay/etc/systemd/system/ados-agent-lite.service
[Unit]
Description=ADOS Drone Agent (lite, Rust runtime)
Documentation=https://github.com/altnautica/ADOSDroneAgent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ados-agent-lite
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Hardening — match the shape used on the Python full agent.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/etc/ados /var/lib/ados /var/log /run/ados

# The lite agent supervises wpa_supplicant + hostapd + dnsmasq when
# the Wi-Fi AP fallback fires. CAP_NET_ADMIN + CAP_SYS_RAWIO let it
# bring up wlan0, configure addresses, and stand up the soft-AP.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_NICE CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
